# etc/profile.d/firstboot.sh
#!/bin/bash

#userdel -r admin
#rm /setup.sh
#rm /etc/profile.d/firstboot.sh
